# cc-tg.service.template
[Unit]
Description=Claude Code Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
ExecStart=__DIR__/start.sh
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Tidak butuh root; tetap minimalisasi privilege.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
